# srcs/influxdb/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:telegraf]
command=telegraf

[program:influxd]
command=influxd
