# srcs/ftps/srcs/supervisord_model.conf
[supervisord]
nodaemon=true

[program:telegraf]
command=telegraf

[program:pure-ftpd]
command=pure-ftpd -j -Y 2 -p 21000:21000 -P __MINIKUBE_IP__
