# srcs/wordpress/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:telegraf]
command=telegraf

[program:lighttpd]
command=lighttpd -D -f /etc/lighttpd/lighttpd.conf
