# srcs/mysql/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:telegraf]
command=telegraf

[program:sh]
command=sh /tmp/start.sh
