# srcs/grafana/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:telegraf]
command=telegraf

[program:grafana-server]
command=grafana-server
