# srcs/nginx/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:telegraf]
command=telegraf

[program:nginx]
command=nginx -g 'daemon off;'
